# systemd-iuvolt.service
[Unit]
Description=Undervolt cpu

[Service]
Type=oneshot
ExecStart=/usr/local/bin/uvolt

[Install]
WantedBy=default.target

