# systemd-iuvolt.service
[Unit]
Description=Undervolt cpu

[Service]
Type=oneshot
ExecStart=/usr/bin/iuvolt

[Install]
WantedBy=default.target

